# systemd/fix_razer_sound.service
[Unit]
Description=Fix Razer Sound

[Service]
ExecStart=/opt/bin/fix_razer_sound.sh

[Install]
WantedBy=multi-user.target
